# trunk-recorder/service/trunk-recorder/run
#!/bin/bash

exec /usr/local/bin/trunk-recorder --config /app/config.json >>/var/log/trunk-recorder.log 2>&1
