# liquidsoap/service/liquidsoap/run
#!/bin/bash

exec liquidsoap /app/west_houston.liq  >>/var/log/icecast.log 2>&1
