# liquidsoap/service/liquidsoap/run
#!/bin/bash

exec liquidsoap /app/stream0.liq  >>/var/log/icecast.log 2>&1
