# icecast/service/icecast/run
#!/bin/bash

exec /sbin/setuser icecast2 icecast2 -c /app/icecast.xml >>/var/log/icecast.log 2>&1
